# clouds/primers/netops_rc
# -- automatically created by BR primer MODIFICATIONS MAY BE LOST
#
# chkconfig: 345 87 13 
# description: netops
#
# processname: netops
# config: ?
# pidfile: /var/run/netops.pid

# !DISCLOSURE! this script is seriously 1/2 baked

# source function library
. /etc/rc.d/init.d/functions

RETVAL=0
PID_FILE=/var/tmp/netops.pid
prog="./init.pl"

# Some functions to make the below more readable
NETOPS_DIR=/home/br/gits/clouds/netops

# BTW -- this is useful
#ps -eo pid,pgid,args | grep pl

test_if_running()
{
    RETVAL=1
    test -f "$PID_FILE" && kill -0 -`cat $PID_FILE` && RETVAL=0
}

start()
{
    test_if_running
    if [ $RETVAL -eq 0 ] ; then
        echo "netops is already running."
        RETVAL=1
    else
        echo -n $"Setting system configuration (DB): "
        /usr/bin/mysql -uroot netops </tmp/br/netops_db_syscfg;
        echo -n $"Starting $prog: "
        ( cd $NETOPS_DIR/br_scripts && /bin/su br -c ./init.pl >../log/init.pl.log 2>../log/init.pl.err) & # -- problem as respawn in init.pl check -t STDOUT
#        ( cd $NETOPS_DIR/br_scripts && /bin/su br -c ./init.pl 2>../log/init.pl.log.err) &
        RETVAL=$?
        echo
    fi
}

stop()
{
    test_if_running
    if [ $RETVAL -eq 0 ] ; then
        echo -n $"Stopping $prog: "
        kill -15 -`cat $PID_FILE` 
        RETVAL=$?
        [ "$RETVAL" = 0 ] && rm -f "$PID_FILE"
        echo
    else
        echo "netops is not running."
        RETVAL=1
    fi
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    restart)
        stop
        start
        ;;
    reload)
        echo "Not implemented";
        RETVAL=1
        ;;
    condrestart)
        echo "Not implemented";
        RETVAL=1
        ;;
    status)
        echo "Not implemented";
        RETVAL=1
        ;;
    *)
        echo $"Usage: $0 {start|stop|restart|reload|condrestart|status}"
        RETVAL=1
esac 
exit $RETVAL

